# supervisord.conf
[supervisord]
nodaemon=true

[program:watcher]
command=python SimpleFlask.py

[eventlistener:httpok]
command=python /usr/local/bin/httpok -p server -t 15 http://localhost/heartbeat
events=TICK_60
